# supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command=cron -f
priority=11

[program:sshd]
command=/usr/sbin/sshd -D -e
autorestart=true
startretries=3
priority=10
